# DEBIAN/prerm
systemctl stop chrony.service
systemctl stop phc2sys.service
systemctl stop ptp4l.service
systemctl stop gpsd.service
systemctl stop timesyncd.service
systemctl disable timesyncd.service
